# packages/bundler-install/templates/systemd.service.template
[Unit]
Description={{DESCRIPTION}}
After=network.target

[Service]
Type=simple
ExecStart={{EXEC_START}}
{{#WORKING_DIR}}
WorkingDirectory={{WORKING_DIR}}
{{/WORKING_DIR}}
Restart={{RESTART}}
RestartSec={{RESTART_DELAY}}
{{#USER}}
User={{USER}}
{{/USER}}
{{#GROUP}}
Group={{GROUP}}
{{/GROUP}}

{{#ENV_VARS}}
Environment="{{KEY}}={{VALUE}}"
{{/ENV_VARS}}

[Install]
WantedBy=multi-user.target
